# apps/visa/visa-helper.service
[Unit]
Description=Bantu Visa Helper Application
After=network.target

[Service]
Type=simple
User=bantu
WorkingDirectory=/home/bantu/dev/bantu-visa-helper
Environment="NODE_ENV=production"
Environment="DATABASE_URL=file:./prisma/prod.db"
ExecStart=/usr/bin/npm start
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
